# com.user.lifetime-reserve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.lifetime-reserve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hongruliu/src/lifetime-reserve/.venv/bin/python</string>
        <string>/Users/hongruliu/src/lifetime-reserve/reserve.py</string>
        <string>--auto</string>
        <string>--wait-until</string>
        <string>09:00:00</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hongruliu/src/lifetime-reserve</string>

    <!-- Run daily at 8:55 AM, login early and book at 9:00 AM sharp -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/hongruliu/src/lifetime-reserve/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hongruliu/src/lifetime-reserve/logs/launchd.log</string>
</dict>
</plist>
